# services/p0ck37/chromedriver/chrome.cron
* * * * * killall -9 --quiet --older-than 1m chrome
