# services/p0ck37/cron/p0ck37.cron
* * * * * find /data -mmin +30 -type f -exec rm -fv {} \;
